# layout/DEBIAN/postinst
#!/bin/sh
launchctl load -w /Library/LaunchDaemons/dev.ayden.ios.lib.sys.aurora.plist || true